# deploy/backup/xfchess-backup.service
[Unit]
Description=XFChess encrypted DB backup (snapshot → age → offsite)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=xfchess
WorkingDirectory=/opt/xfchess
EnvironmentFile=/opt/xfchess/.env
ExecStart=/opt/xfchess/deploy/backup/backup-db.sh

# Hardening (this unit only reads the DBs and writes to the backups workdir).
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/xfchess/data
ProtectHome=true
PrivateTmp=true
